# supervisord_programs/bioc-webstats.service
# TODO Hard coded for standalone deployment
# usage:
#     sudo systemctl enable bioc-webstats.service
#     sudo systemctl start bioc-webstats.service
[Unit]
Description=bioc-webstats backend web server
After=network.target

[Service]
ExecStart=/home/ubuntu/bio-web-stats/supervisord_programs/start_gunicorn.sh
Restart=always
User=ubuntu
Group=ubuntu
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target